# recipes-support/edgeos-mounts/edgeos-mounts/greengrass-ggc-deployment.mount
[Unit]
Description=greengrass lambda directory bind mount
Requires=mnt-data.mount edgeos-data-mnt-init.service
After=mnt-data.mount edgeos-data-mnt-init.service

[Mount]
What=/mnt/data/root-overlay/greengrass/ggc/deployment
Where=/greengrass/ggc/deployment
Type=none
Options=bind

[Install]
WantedBy=edgeos-bind.target